# services/json2mqtt/run
#!/bin/sh

# -u unbuffered output
exec ../../venv/bin/python -u ../../json2mqtt/run.py 2>&1
