# services/json2mqtt/run
#!/bin/sh

# -u unbuffered output
exec ../../../feinstaub-json-mqtt-bridge-ve/bin/python -u ../../json2mqtt/run.py 2>&1
